# deploy/systemd/nostrig-serve@.service
[Unit]
Description=nostrig task-fabric recipient (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
Environment=NOSTRIG_ENV=production
Environment=NOSTR_RELAY=wss://relay.sharegap.net
Environment=NOSTRIG_ACL_FILE=%d/nostrig_acl
Environment=NOSTRIG_SIGNER_CLIENT_SECRET_KEY_FILE=%d/nostrig_signer_client_secret_key
EnvironmentFile=/etc/nostrig/%i.env
LoadCredential=nostrig_acl:/etc/nostrig/%i.acl.json
LoadCredential=nostrig_signer_client_secret_key:/etc/nostrig/%i.signer-client-secret
RuntimeDirectory=nostrig-%i
RuntimeDirectoryMode=0700
StateDirectory=nostrig-%i
StateDirectoryMode=0700
UMask=0077
ExecStart=/usr/local/bin/nostrig serve --observability-addr 127.0.0.1:8080 --outbox-path /var/lib/nostrig-%i/outbox.json --instance-lock /var/lib/nostrig-%i/instance.lock
ExecStartPost=/usr/bin/curl --fail --silent --show-error --retry 120 --retry-delay 1 --retry-connrefused --retry-all-errors --connect-timeout 2 --max-time 180 http://127.0.0.1:8080/readyz
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
